# Semester6/SOS/C07/iptables.service
[Unit]
Description=Load iptables
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/init_firewall.sh

[Install]
WantedBy=multi-user.target
